# docker/s6/services.d/apache/run
#!/usr/bin/execlineb -P

with-contenv
/usr/sbin/apachectl -DFOREGROUND